# vnv_manager/supervisord.conf
[supervisord]
user=root

[program:nginx]
command=/usr/sbin/nginx
autorestart=false
#stdout_logfile=/dev/stdout
#stdout_logfile_maxbytes=0
#stderr_logfile=/dev/stderr
#stderr_logfile_maxbytes=0

[program:uwsgi]
command=/usr/sbin/uwsgi --ini /etc/uwsgi/uwsgi.ini --die-on-term
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[program:migrate_db]
command = /opt/Monitoring/run.sh
startsecs = 0
autorestart = false
startretries = 1

[program:osm_plugin]
command = /usr/bin/python /opt/Monitoring/manage.py osm_plugin --traceback --verbosity=2
startsecs = 0
autorestart = false
startretries = 1
