# pushgateway/supervisord.conf
[supervisord]
nodaemon=true

[program:pushgateway]
command = /opt/Monitoring/pushgateway/pushgateway
startsecs = 0
autorestart = true
startretries = 1